# FlightTracker3DApp/readsb.service
[Unit]
Description=Readsb ADS-B Decoder
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/readsb
ExecStart=/home/pi/readsb/readsb --device-type rtlsdr --gain auto --interactive --net --net-sbs-port=30003

# Auto-restart configuration
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=readsb

# Environment variables
Environment="DOTNET_CLI_HOME=/tmp"

[Install]
WantedBy=multi-user.target
